# package/base-files/files/etc/crontabs/root
*/2 * * * * i2cset -y 1 0x63 0x21 0x22 0x55 0x00 0x33 0x44 0x55 0x66 i



# crontab must end with the last line as space or a comment

